# crontab
# Pobieraj oferty codziennie o 06:00 w nocy i o 14:00 po południu
0 6,14 * * * curl -X POST http://127.0.0.1:8000/zbieraj/ >> /home/admin/job_aggregator/cron_log.txt 2>&1
